# deploy/hmail-api.service
[Unit]
Description=PMail+ API (Prohost Cloud)
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/hostnet-panel/apps/hmail-api
EnvironmentFile=/var/www/hostnet-panel/.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
